# com.ccon.streamlit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ccon.streamlit</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/Samsung USB/cardnews_3/start_streamlit.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/Samsung USB/cardnews_3</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Volumes/Samsung USB/cardnews_3/logs/streamlit_service.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/Samsung USB/cardnews_3/logs/streamlit_service_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>

